# my-claude-bot@.service
[Unit]
Description=My Claude Bot — Telegram AI Assistants
# Использование: sudo systemctl enable my-claude-bot@YOUR_USERNAME
# Пример: sudo systemctl enable my-claude-bot@ubuntu
After=network.target

[Service]
Type=simple
# Замените на вашего пользователя и путь к проекту
User=%i
WorkingDirectory=/home/%i/my-claude-bot
ExecStart=/usr/bin/python3 -m src.main
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Безопасность
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/%i/my-claude-bot/agents

[Install]
WantedBy=multi-user.target
